# rootfs/etc/supervisor/conf.d/aptly.conf
[program:aptly-api]
command=/usr/bin/aptly api serve -listen=:8080 -no-lock
environment=GIN_MODE=release
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0